# .docker/prod.supervisord.conf
[program:laravel-server]
command=php /app/artisan octane:frankenphp --host=0.0.0.0 --port=80 --admin-port=2019 --quiet
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0

[program:laravel-queue]
command=php /app/artisan queue:work --timeout=3600 --quiet
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0
